# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-java21-add-package/run
#!/usr/bin/with-contenv bash

echo "**** Adding java21 packages to install list ****"
echo "openjdk-21-jdk" >> /mod-repo-packages-to-install.list
