# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-java8-add-package/run
#!/usr/bin/with-contenv bash

echo "**** Adding java8 packages to install list ****"
echo "openjdk-8-jdk" >> /mod-repo-packages-to-install.list
